# deploy/systemd/audiolad-maintenance.timer
[Unit]
Description=Daily Audiolad disk hygiene (off-peak UTC)

[Timer]
OnCalendar=*-*-* 04:30:00
Persistent=true
RandomizedDelaySec=900

[Install]
WantedBy=timers.target
